# pacman/copy.list
etc/bash
etc/xonsh
etc/default/lxc-net
etc/default/grub
etc/hdparm.conf
etc/modprobe.d
etc/modules-load.d
etc/dracut.conf.d
etc/refractainstaller.conf
etc/refractasnapshot.conf
etc/skel
etc/sddm.conf.d
etc/ssh/sshd_config.d
etc/sysctl.d
etc/systemd/logind.conf
etc/systemd/resolved.conf.d
etc/sudoers
etc/tor
usr/local/bin
usr/share/lxc/templates/lxc-salt_tarball
usr/share/X11/xkb
etc/systemd/system/relabel-selinux.service
